# system/init/bot.conf
description "Bot"

start on net-device-up
stop on runlevel [016]

respawn
respawn limit 5 30

setuid bot

script

    exec /usr/local/bin/php /usr/local/bin/bot/bot.php

end script
